# root/etc/services.d/thumbnails-monitor/run
#!/usr/bin/execlineb -P

# Start daemon
/thumbnails-monitor.sh
