# root/etc/services.d/thumbnails-monitor/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid thumbnailer /thumbnails-monitor.sh
