# root/etc/services.d/thumbnails-monitor/run
#!/usr/bin/execlineb -P

s6-setuidgid minidlna

# Start daemon
/thumbnails-monitor.sh
