# root/etc/services.d/thumbnails-monitor/run
#!/usr/bin/with-contenv bash

exec /thumbnails-monitor.sh
