# root/etc/services.d/minidlna/run
#!/usr/bin/with-contenv bash

# Start daemon
exec exec s6-setuidgid minidlna /usr/sbin/minidlnad -S