# root/etc/services.d/minidlna/run
#!/usr/bin/with-contenv bash

# Start daemon
exec /usr/sbin/minidlnad -S -u minidlna