# root/etc/services.d/minidlna/run
#!/usr/bin/execlineb -P

# Start daemon
/usr/sbin/minidlnad -S -r 